# deployconfig/backupdb.crontab
*/1 * * * * python {remote_website_dir}/tool/backupdb.py
